# charts/slurm/scripts/profile.d/utils.sh
#!/bin/bash

# set path to some utilities that we manage

export PATH="/etc/slurm-cm-exe:/home/common/git-sync/k8s/k8s.git/bin:$PATH"
